# installers/linux/DEBIAN/postinst
#!/bin/sh
set -e

# Create necessary directories
mkdir -p /etc/sssonector/certs
mkdir -p /var/log/sssonector

# Set up systemd service
if [ -d /lib/systemd/system ]; then
    cat > /lib/systemd/system/sssonector.service << EOF
[Unit]
Description=SSSonector SSL Tunnel Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sssonector -config /etc/sssonector/config.yaml
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target
EOF

    # Reload systemd
    systemctl daemon-reload

    # Enable service but don't start it (user needs to configure first)
    systemctl enable sssonector.service
fi

# Set permissions
chmod 755 /etc/sssonector
chmod 700 /etc/sssonector/certs
chmod 755 /var/log/sssonector

# Create default config if it doesn't exist
if [ ! -f /etc/sssonector/config.yaml ]; then
    cp /usr/share/sssonector/config/server.yaml /etc/sssonector/config.yaml
fi

echo "SSSonector installation completed."
echo "Please configure your certificates in /etc/sssonector/certs/"
echo "and update /etc/sssonector/config.yaml before starting the service."
