# test/qa_scripts/sssonector.service
[Unit]
Description=SSSonector Secure Tunnel Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
RuntimeDirectory=sssonector
RuntimeDirectoryMode=0770
ExecStartPre=-/usr/sbin/ip link delete tun0
ExecStartPre=/usr/sbin/ip tuntap add dev tun0 mode tun user root
ExecStartPre=/usr/sbin/ip link set tun0 up
ExecStart=/usr/local/bin/sssonector -config /etc/sssonector/config.yaml
ExecStopPost=-/usr/sbin/ip link delete tun0
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal

# Security settings
ProtectSystem=no
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=no
DevicePolicy=closed
DeviceAllow=/dev/net/tun rw
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
